# init/runit/mvm-instance/finish
#!/bin/sh
# Stop instance when the service is taken down.
exec "@MVM_ROOT@/mvm" stop "@INSTANCE@"
