# init/runit/mvm-argus/finish
#!/bin/sh
# runit finish for Argus (flush on down).
exec "@MVM_ROOT@/mvm" argus flush
